# crontab
* * * * * docker ps >> /var/log/cron/cron.log 2>&1
*/2 * * * * echo "Hello ${ENV_NAME}" >> /var/log/cron/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
